# ros1/image/app/conf.d/start.conf
[program:start]
command=/app/bin/start.sh
autorestart=true
stdout_logfile=/app/logs/start.log
redirect_stderr=true